# arch-update-popup.service
[Unit]
Description=Arch Linux update tray notifier
After=graphical-session.target network-online.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/arch-update-popup
Restart=on-failure
RestartSec=30

[Install]
WantedBy=graphical-session.target
