# arch-update-popup.timer
[Unit]
Description=Start Arch Linux update tray notifier after login
Requires=arch-update-popup.service

[Timer]
# Start 5 minutes after the graphical session starts
OnStartupSec=5min
# App handles its own periodic checks internally

[Install]
WantedBy=timers.target
